# deploy/com.lark-codex-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lark-codex-gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__GATEWAY_RUNNER__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__GATEWAY_ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__GATEWAY_ROOT__/gateway.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__GATEWAY_ROOT__/gateway.stderr.log</string>
</dict>
</plist>
